# target/linux/ipq60xx/base-files/etc/hotplug.d/firmware/10-ath11k-bdf
#!/bin/sh

[ -e /lib/firmware/$FIRMWARE ] && exit 0

case "$FIRMWARE" in
"ath11k/IPQ6018/hw1.0/board-2.bin")
	BOARD_DATA="ath11k/IPQ6018/hw1.0/$(board_name)-board-2.bin"
	if [ -e "/lib/firmware/$BOARD_DATA" ]; then
		ln -s "/lib/firmware/$BOARD_DATA" "/lib/firmware/$FIRMWARE"
	fi
esac
